# Ansible-Jumbo/files/web-cs-prod/etc/supervisor/conf.d/cs-ee.conf
[program:cs-ee]
command=java -Xms2048M -Xmx4096M -jar /var/www/CONTENTSERV/Jumbo_Prod/admin/core/extensions/exportstaging/java/dist/ExportExecutor-16.0.jar
directory=/var/www/CONTENTSERV/Jumbo_Prod/admin/core/extensions/exportstaging/java/dist/
priority=30
autostart=true
user=www-data
startsecs=20
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/cs-ee.log
stderr_logfile=/var/log/supervisor/cs-ee.err

