# CS-16.1-Ubuntu-Xenial-Installer/additions/etc/supervisor/conf.d/cs-ee.conf
[program:cs-ee]
command=java -Xms1024M -Xmx1024M -jar /var/lib/contentserv-cs/admin/core/extensions/exportstaging/java/dist/ExportExecutor-16.1.jar
directory=/var/lib/contentserv-cs/admin/core/extensions/exportstaging/java/dist/
priority=30
autostart=false
user=www-data
startsecs=20
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/cs-ee.log
stderr_logfile=/var/log/supervisor/cs-ee.err

